# PNPros/ROS_bridge/pnp_ros/qtcreator-build/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/jnitsch/ros_hydro/pnp_ws/src/petri_net/pnp_ros/qtcreator-build/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
